# UserData/IB Support/Simulator Devices/0BC9BF9B-0AD3-4304-A8BA-0778331A826C/data/Library/LaunchDaemons/com.apple.sketchboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>EmbeddedHomeScreen</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.sketchboard</string>
	<key>MachServices</key>
	<dict>
		<key>PurpleSystemAppPort</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.UIKit.KeyboardManagement.hosted</key>
		<true/>
		<key>com.apple.UIKit.statusbarserver</key>
		<true/>
		<key>com.apple.bulletinboard.dataproviderconnection</key>
		<true/>
		<key>com.apple.bulletinboard.observerconnection</key>
		<true/>
		<key>com.apple.bulletinboard.serverconduitconnection</key>
		<true/>
		<key>com.apple.bulletinboard.settingsconnection</key>
		<true/>
		<key>com.apple.bulletinboard.systemstateconnection</key>
		<true/>
		<key>com.apple.bulletinboard.utilitiesconnection</key>
		<true/>
		<key>com.apple.frontboard.systemappservices</key>
		<true/>
		<key>com.apple.frontboard.workspace</key>
		<true/>
		<key>com.apple.sketchboard.services.guest</key>
		<true/>
		<key>com.apple.sketchboard.services.host</key>
		<true/>
		<key>com.apple.usernotification.notificationregistrar</key>
		<true/>
		<key>com.apple.usernotification.notificationscheduler</key>
		<true/>
		<key>com.apple.usernotifications.remotenotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationservice</key>
		<true/>
		<key>com.apple.usernotifications.usernotificationsettingsservice</key>
		<true/>
	</dict>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>Program</key>
	<string>/Applications/Xcode-beta.app/Contents/Developer/Platforms/iPhoneOS.platform/Library/Developer/CoreSimulator/Applications/SketchBoard.app/SketchBoard</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/SketchBoard.app/SketchBoard</string>
		<string>system_app_start</string>
	</array>
	<key>RootedSimulatorPath</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>10</integer>
</dict>
</plist>
